# configs/systemd/user/sommelier-x@0.service.d/cros-sommelier-x-override.conf
[Service]
Environment="SOMMELIER_FRAME_COLOR=#1C2D33"
Environment="SOMMELIER_ACCELERATORS=Super_L,<Alt>bracketleft,<Alt>bracketright,<Alt>equal"
