# configs/systemd/user/dconf-backup.path
[Path]
PathModified=%h/.config/dconf

[Install]
WantedBy=default.target