# configs/systemd/user/etc-hosts-watcher.path
[Path]
PathModified=/etc/hosts

[Install]
WantedBy=multi-user.target