# configs/systemd/user/rclone-jobs.service
[Unit]
Description="Run rclone jobs listed in beleyenv's config.json."

[Service]
Type=oneshot
ExecStart="%h/.beleyenv/beleyenv/bin/runRcloneJobs"